# config/launchd/com.local.kokoro.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  NOTE: macOS plists do not support $HOME or ~.
  Replace __HOME__ with your actual home directory before installing:
    sed "s|__HOME__|$HOME|g" com.local.kokoro.plist > ~/Library/LaunchAgents/com.local.kokoro.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.kokoro</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/.kokoro-fastapi/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/.kokoro-fastapi</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/kokoro.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/kokoro.stderr.log</string>
</dict>
</plist>
